# deploy/systemd/wifi-connect.service
[Unit]
Description=WiFi onboarding captive portal (balena wifi-connect)
After=NetworkManager.service
Wants=NetworkManager.service
Before=embedded-music.service

[Service]
Type=oneshot
RemainAfterExit=no
# wifi-connect exits immediately if already connected to a known network;
# otherwise it raises the AP + portal and blocks until the user supplies creds.
Environment=PORTAL_SSID=Embedded-Music-Setup
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/sbin/wifi-connect --portal-ssid ${PORTAL_SSID} --ui-directory /usr/local/share/wifi-connect/ui
# No timeout: an un-provisioned box should wait for setup rather than fail boot.
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
